# etc/udev/rules.d/z98CanonEOS400D-stuart.rules
# Stuart
ACTION=="add", ATTRS{idVendor}=="04a9", ATTRS{idProduct}=="3110", RUN+="/usr/local/sbin/canon_download_term $env{idVendor}"
